# systemd/gpu-stats-poll.service
[Unit]
Description=Poll GPU stats from mining server
After=network.target

[Service]
Type=oneshot
ExecStart=/home/bowen/ha-smartshift/scripts/poll_gpu_stats.sh
